# config/green_controller_service.service.example
[Unit]
Description=Green Controller Service

[Service]
ExecStart=make -C /home/admin/code/green_controller_module/ run

Type=notify

[Install]
WantedBy=default.target